# templates/ai-proxy-firewall.service
[Unit]
Description=Allow AI platform Docker containers to reach host Privoxy proxy
After=ufw.service docker.service privoxy.service
Wants=ufw.service docker.service privoxy.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/ai-proxy-firewall.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

